# init/god.conf
# Upstart configuation
description "god daemon"

start on runlevel [2]
stop on runlevel [016]

exec /usr/bin/god -c /etc/god/master.conf -D

respawn
